# fan_remote.socket
[Socket]
ListenStream=23734

[Install]
WantedBy=sockets.target
